# systemd/prometheus-smart-collector.service
[Unit]
Description=Prometheus SMART collector

[Service]
Type=oneshot
ExecStart=prometheus-smart-collector

[Install]
WantedBy=multi-user.target
